# scripts/systemd/openclaw-auth-monitor.service
[Unit]
Description=OpenClaw Auth Expiry Monitor
After=network.target

[Service]
Type=oneshot
ExecStart=/home/admin/openclaw/scripts/auth-monitor.sh
# Configure notification channels via environment
Environment=WARN_HOURS=2
# Environment=NOTIFY_PHONE=+1234567890
# Environment=NOTIFY_NTFY=openclaw-alerts

[Install]
WantedBy=default.target
