# scripts/systemd/openclaw-auth-monitor.timer
[Unit]
Description=Check OpenClaw auth expiry every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
